# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── nginx ────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stderr_logfile=/var/log/supervisor/nginx.stderr.log

; ── gunicorn ─────────────────────────────────────────────────────────────────
[program:gunicorn]
command=gunicorn
    --workers %(ENV_GUNICORN_WORKERS)s
    --bind 127.0.0.1:%(ENV_GUNICORN_PORT)s
    --timeout 3600
    --worker-class sync
    --log-level info
    --access-logfile -
    --error-logfile -
    "app:create_app()"
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/gunicorn.stdout.log
stderr_logfile=/var/log/supervisor/gunicorn.stderr.log

; ── expiry sweep (runs every hour via a simple loop) ─────────────────────────
[program:expiry]
command=bash -c "while true; do python -m tasks.expiry; sleep 3600; done"
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/expiry.stdout.log
stderr_logfile=/var/log/supervisor/expiry.stderr.log
